# scripts/com.cratedigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cratedigger</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/elvista/Sites/SpotDownload/backend/venv/bin/uvicorn</string>
    <string>main:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>5174</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/elvista/Sites/SpotDownload/backend</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/elvista/Library/Logs/cratedigger/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/elvista/Library/Logs/cratedigger/stderr.log</string>
</dict>
</plist>
